# recipes-system/service-options/files/service-options-pendrive.service.in
[Unit]
Description="Service options"
ConditionFileNotEmpty=/mnt/pendrive/boot/service.conf
Requires=mnt-pendrive.mount
Wants=multi-user.target service-options@mnt-pendrive-boot-service.conf.service
Before=service-options@mnt-pendrive-boot-service.conf.service
After=multi-user.target mnt-pendrive.mount

[Service]
ExecStart=@bindir@/true

[Install]
WantedBy=multi-user.target
