# recipes-system/service-options/files/service-options.service.in
[Unit]
Description="Service options"
Requires=multi-user.target
After=multi-user.target

[Service]
Type=oneshot
ExecStart=@sbindir@/service-options /rootfs.dev/boot/service.conf

[Install]
WantedBy=multi-user.target
